# templates/tomcat.service.erb
[Unit]
Description=Tomcat Service

[Service]
Type=forking
ExecStart=/etc/init.d/<%= @scriptName %>

[Install]
WantedBy=multi-user.target
